# base/mkosi.extra/usr/lib/systemd/system/systemd-growfs-root.service.d/50-pocketfed-ostree.conf
[Unit]
ConditionKernelCommandLine=ostree
After=ostree-remount.service
