# base/mkosi.extra/usr/lib/systemd/system/tqftpserv.service.d/blob-wrangler.conf
[Unit]
Requires=blob-wrangler.service
After=blob-wrangler.service
