# devices/google-crosshatch/pocketfed-usb-developer-mode.service
[Unit]
Description=Activate the PocketFed Crosshatch USB developer link
Requires=usb-signaller.service
After=usb-signaller.service
PartOf=usb-signaller.service

[Service]
Type=oneshot
ExecStart=/usr/libexec/pocketfed-usb-signaller-activate
RemainAfterExit=yes
Restart=on-failure
RestartSec=2s
TimeoutStartSec=90s
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
