# base/mkosi.extra/usr/lib/systemd/system/blob-wrangler.service.d/post-modprobe.conf
[Service]
ExecStartPost=/usr/bin/bash /usr/libexec/pocketfed-blob-wrangler-modprobe
