# devices/google-crosshatch/blob-wrangler.service.d/10-pocketfed-crosshatch-slot-fallback.conf
[Service]
ExecStart=
ExecStart=/usr/libexec/pocketfed-crosshatch-blob-wrangler
