# devices/google-sargo/usb-signaller.service.d/10-pocketfed-sargo.conf
[Service]
# Ignore mutable /etc policy left by earlier images and select the immutable
# sargo management-link policy shipped with this deployment.
EnvironmentFile=
EnvironmentFile=/usr/lib/pocketfed/usb-signaller-sargo.conf
